# deploy/findthisplace-update.service
[Unit]
Description=Update findthisplace binary

[Service]
Type=oneshot
ExecStart=/var/www/findthisplace.eu/update.sh
