# deploy/findthisplace-update.path
[Unit]
Description=Watch for findthisplace update trigger

[Path]
PathExists=/run/findthisplace-update/trigger

[Install]
WantedBy=multi-user.target
