# deploy/findthisplace-update.path
[Unit]
Description=Watch for findthisplace update trigger

[Path]
PathExists=/tmp/findthisplace-update-trigger

[Install]
WantedBy=multi-user.target
